# stiebel-control.service
[Unit]
Description=Stiebel Eltron Heat Pump Control
After=network.target

[Service]
User=your_user
WorkingDirectory=/path/to/stiebel-control-python
ExecStart=/usr/bin/python3 -m stiebel_control --config service_config.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
